# limon.service
[Unit]
Description=LiMON AI Agent - System startup service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/limon
ExecStart=/usr/bin/node /opt/limon/index.js
Restart=on-failure
RestartSec=10s

# Environment
Environment="NODE_ENV=production"
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
